# .config/systemd/user/tldr-cache-update.service
[Unit]
Description=Periodically update the tldr cache

[Service]
Type=oneshot
ExecStart=/usr/bin/tldr --update
